# deploy/raspberry-pi/anchor-mics-kiosk.service
[Unit]
Description=Anchor Mics Chromium kiosk
After=graphical.target anchor-mics.service
Wants=anchor-mics.service

[Service]
Type=simple
User=cci
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/cci/.Xauthority
Environment=ANCHOR_MICS_KIOSK_URL=http://127.0.0.1:8010/display
WorkingDirectory=/opt/anchor-mics
ExecStart=/bin/bash /opt/anchor-mics/deploy/raspberry-pi/start-kiosk.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
